# docs/systemd/paycrivo-web.service
# PayCrivo web (TanStack Start SSR) — systemd unit
# Install to: /etc/systemd/system/paycrivo-web.service
# Then: sudo systemctl daemon-reload && sudo systemctl enable --now paycrivo-web
#
# This runs PayCrivo's hardened Node SSR launcher. The launcher serves
# .output/public assets and dispatches application routes directly to the
# generated TanStack Start SSR service at .output/server/_ssr/ssr.mjs.
#
# The TanStack app lives at the project root (/var/www/paycrivo.com). The
# separate Express API (server/) runs as paycrivo-api on port 4100.

[Unit]
Description=PayCrivo Web (TanStack Start SSR)
After=network.target

[Service]
Type=simple
WorkingDirectory=/var/www/paycrivo.com
# SSR server port (reverse-proxied by Apache). Change if 4000 is taken.
Environment=PORT=4000
Environment=HOST=127.0.0.1
Environment=NODE_ENV=production
ExecStart=/usr/bin/node scripts/start-web.mjs
Restart=always
RestartSec=5
User=www-data

[Install]
WantedBy=multi-user.target
